# webdav/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: WebDAV
# Runs caddy
# ==============================================================================

bashio::log.info "Starting Caddy..."

exec caddy run --config /etc/caddy/Caddyfile --adapter caddyfile