# webdav/rootfs/etc/s6-overlay/s6-rc.d/rclone/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: WebDAV
# WebDAV init script, runs before any other service
# This scripts runs a WebDAV server using rclone
# ==============================================================================
declare username
declare password
declare certfile
declare keyfile

bashio::config.require.ssl

bashio::log.info "Initializing rclone.."
username=$(bashio::config 'username')
password=$(bashio::config 'password')
certfile=$(bashio::config 'certfile')
keyfile=$(bashio::config 'keyfile')

# rclone SSL options
rclone_ssl_options=""
if bashio::config.true 'ssl'; then
	rclone_ssl_options=" --cert /ssl/${certfile} --key /ssl/${keyfile} "
fi

rclone serve webdav --user ${username} \
                    --pass ${password} \
                    --addr :8765 \
                    ${rclone_ssl_options} \
                    /data