# webdav/rootfs/etc/s6-overlay/s6-rc.d/init-caddy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: WebDAV
# WebDAV init script, runs before any other service
# This scripts generates password hash and prepares the caddy configuration file
# ==============================================================================
declare username
declare password
declare certfile
declare keyfile

bashio::config.require.ssl

bashio::log.info "Initializing caddy configuration"
username=$(bashio::config 'username')
password=$(bashio::config 'password')
certfile=$(bashio::config 'certfile')
keyfile=$(bashio::config 'keyfile')

# calculate password hash
hashed_password=$(echo ${password} | caddy hash-password)

if bashio::config.true 'ssl'; then
	sed -e "s/WEBDAV_USER/${username}/g" \
		-e "s@WEBDAV_PASSWORD@${hashed_password}@g" \
		-e "s@CERTFILE@${certfile}@g" \
		-e "s@KEYFILE@${keyfile}@g" \
		/etc/caddy/Caddyfile.tmpl > /etc/caddy/Caddyfile
else 
	sed -e "s/WEBDAV_USER/${username}/g" \
		-e "s@WEBDAV_PASSWORD@${hashed_password}@g" \
		/etc/caddy/Caddyfile.nossl.tmpl > /etc/caddy/Caddyfile
fi